# docroot/pubkey-watcher.service
[Unit]
Description=Watch /var/www/html/pubkey and update users.list
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/html
ExecStart=/var/www/html/pubkey-watcher.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
